# .crontab
*/5 * * * * /usr/local/bin/update_ip
0-59 * * * * /usr/local/bin/ddns_update
